# alsa_in
#!/bin/sh
### BEGIN INIT INFO
# chkconfig: - 95 15
# Provides:       Alsa-in Server
# Required-Start:
# Required-Stop:
# Default-Start:  3 4 5
# Default-Stop:
# Description:    Alsa-in Server
### END INIT INFO

# Source function library.
if [ -f /etc/init.d/functions ] ; then
  . /etc/init.d/functions
elif [ -f /etc/rc.d/init.d/functions ] ; then
  . /etc/rc.d/init.d/functions
else
  exit 1
fi


cd /tmp

# Intamixx Config.
. /intamixx/intamixxdj.conf

deck1_cardnum=`awk '{i1=$0; getline ; i2=$0 ; x= i1 i2 ; print x}' /proc/asound/cards | grep -i $xwax_deck1 | awk -F' ' '{print $1}'`

BIN_DIR="/usr/bin"
ARGS="-v -j $slave_card_in -d hw:$deck1_cardnum -c 2 -p 48000 -p 256 -n 3"
USER=intamixx

# List process here
Process="alsa_in"
DisplayName="Alsa-In"

# Source function library.
. /etc/rc.d/init.d/functions

# Builtin functions need to add terminating newline
function psuccess
{
        success; echo ""
}
function pfailure
{
        failure; echo ""
}
function pwarning
{
        warning; echo ""
}

function start
{

        for x in $Process; do
                if [[ -x $BIN_DIR/$x ]]; then
                        if [ -n "`/sbin/pidof $x`" ];
                                then
                                        echo -n "$DisplayName is already running."
                                else
                                        echo "Starting $DisplayName..."
					#if [ -e "$DEV" ]
					#then
                                        #	$BIN_DIR/$x > $DEV &
					#else
                                        /usr/bin/xterm -geometry 60x10 -fg yellow -bg black -T "Intamixx Alsa-In" -e "$BIN_DIR/$x $ARGS" &
                                        #$BIN_DIR/$x $ARGS &
					#fi
			fi

			sleep 5

			if [ -n "`/sbin/pidof $x`" ];
				then
					psuccess
					retcode=0
				else
					pfailure
					retcode=1
                        fi
                fi
        done
}

function stop
{

        for x in $Process; do
                if [[ -x $BIN_DIR/$x ]]; then
                        if [ -n "`/sbin/pidof $x`" ];
                                then
                                        echo "Stopping $DisplayName..."
                                        killall $x
                                else
                                        echo "$DisplayName is already stopped."
                        fi

			sleep 1

			if [ -n "`/sbin/pidof $x`" ];
				then
					pfailure
					retcode=1
				else
					psuccess
					retcode=0
                        fi
                fi
        done
}

function status
{
        for x in $Process; do
                if [[ -x $BIN_DIR/$x ]]; then
                        if [  -n "`/sbin/pidof $x`" ];
                                then
                                        echo "$DisplayName is running"
					psuccess
					retcode=0
                                else
                                        echo "$DisplayName is stopped"
					pfailure
					retcode=1
                        fi
                fi
        done
}

function usage
{
        echo "Usage: $0 (start|stop|restart|status)"
}

case $1 in
        start)
                start
                ;;
        stop)
                stop
                ;;
        restart)
                stop && start
                ;;
        status)
                status
                ;;
        *)
                usage
                retcode=1
esac
exit $retcode
 
